# docs/配置文件模板/supervisord.fqnext.example.conf
[supervisord]
; FreshQuant 当前宿主机 Python 口径：项目根目录 .venv
; 如仓库路径不同，请同步调整下面的 PATH、PYTHONPATH 和 command
logfile=D:/fqdata/log/supervisord_fqnext.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=error
pidfile=D:/fqdata/supervisord_fqnext.pid
identifier=fqnext

[program-default]
directory=D:/fqpack/freshquant-2026.2.23
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
autorestart=true
startretries=10
envFiles=D:/fqpack/config/envs.conf
environment=PATH=D:/fqpack/freshquant-2026.2.23/.venv;D:/fqpack/freshquant-2026.2.23/.venv/Scripts;C:/Windows/System32,PYTHONPATH=D:/fqpack/freshquant-2026.2.23;D:/fqpack/freshquant-2026.2.23/morningglory/fqxtrade

[inet_http_server]
port=127.0.0.1:10011

[program:fqnext_realtime_xtdata_producer]
command=D:/fqpack/freshquant-2026.2.23/.venv/Scripts/python.exe -m freshquant.market_data.xtdata.market_producer
directory=D:/fqpack/freshquant-2026.2.23
stdout_logfile=D:/fqdata/log/fqnext_realtime_xtdata_producer.log
stderr_logfile=D:/fqdata/log/fqnext_realtime_xtdata_producer_err.log
autostart=true
autorestart=true
startsecs=5

[program:fqnext_realtime_xtdata_consumer]
command=D:/fqpack/freshquant-2026.2.23/.venv/Scripts/python.exe -m freshquant.market_data.xtdata.strategy_consumer --prewarm --max-bars 20000
directory=D:/fqpack/freshquant-2026.2.23
stdout_logfile=D:/fqdata/log/fqnext_realtime_xtdata_consumer.log
stderr_logfile=D:/fqdata/log/fqnext_realtime_xtdata_consumer_err.log
autostart=true
autorestart=true
startsecs=5

[program:fqnext_xtquant_broker]
command=D:/fqpack/freshquant-2026.2.23/.venv/Scripts/python.exe -m fqxtrade.xtquant.broker
directory=D:/fqpack/freshquant-2026.2.23
stdout_logfile=D:/fqdata/log/fqnext_xtquant_broker.log
stderr_logfile=D:/fqdata/log/fqnext_xtquant_broker_err.log
autostart=true
autorestart=true
startsecs=5

[program:fqnext_credit_subjects_worker]
command=D:/fqpack/freshquant-2026.2.23/.venv/Scripts/python.exe -m freshquant.order_management.credit_subjects.worker
directory=D:/fqpack/freshquant-2026.2.23
stdout_logfile=D:/fqdata/log/fqnext_credit_subjects_worker.log
stderr_logfile=D:/fqdata/log/fqnext_credit_subjects_worker_err.log
autostart=true
autorestart=true
startsecs=5

[program:fqnext_xtdata_adj_refresh_worker]
command=D:/fqpack/freshquant-2026.2.23/.venv/Scripts/python.exe -m freshquant.market_data.xtdata.adj_refresh_worker
directory=D:/fqpack/freshquant-2026.2.23
stdout_logfile=D:/fqdata/log/fqnext_xtdata_adj_refresh_worker.log
stderr_logfile=D:/fqdata/log/fqnext_xtdata_adj_refresh_worker_err.log
autostart=true
autorestart=true
startsecs=5

[group:fqnext_reference_data]
programs=fqnext_credit_subjects_worker,fqnext_xtdata_adj_refresh_worker

[supervisorctl]
serverurl=http://127.0.0.1:10011
